# ops/systemd/labot-web.service
[Unit]
Description=LaBot Web Application
After=network.target postgresql.service

[Service]
Type=notify
User=labot
Group=labot
WorkingDirectory=/opt/labot/web
EnvironmentFile=/etc/labot/web.env
ExecStart=/usr/bin/dotnet /opt/labot/web/LaBot.Web.dll
Restart=on-failure
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=labot-web
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
