# ops/systemd/labot-worker.service
[Unit]
Description=LaBot Worker Service (Bot Engine)
After=network.target postgresql.service

[Service]
Type=simple
User=labot
Group=labot
WorkingDirectory=/opt/labot/worker
EnvironmentFile=/etc/labot/worker.env
ExecStart=/usr/bin/dotnet /opt/labot/worker/LaBot.Worker.dll
Restart=on-failure
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=labot-worker
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
